# templates/zeppelin.upstart.conf.j2
description "zeppelin"

start on (local-filesystems and net-device-up IFACE!=lo)
stop on shutdown

respawn
respawn limit 7 5

chdir /opt/zeppelin
exec bin/zeppelin-daemon.sh upstart
